# imageroot/systemd/user/flexisip.service
[Unit]
Description=NethVoice flexisip.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
EnvironmentFile=%S/state/environment
Restart=always
WorkingDirectory=%S/state
TimeoutStopSec=70
ExecStartPre=/bin/rm -f %t/flexisip.pid %t/flexisip.ctr-id
ExecStartPre=-runagent install-certificate flexisip
ExecStart=/usr/bin/podman run \
  --conmon-pidfile=%t/flexisip.pid \
  --cidfile=%t/flexisip.ctr-id \
  --replace --name=flexisip \
  --detach \
  --network=host \
  --volume=asterisk:/etc/asterisk:z \
  --volume=flexisip-cert:/etc/flexisip/tls:z \
  --env NETHVOICE_HOST \
  --env FLEXISIP_* \
  --env REDIS_FLEXISIP_PORT \
  --env APACHE_PORT \
  --env SUBSCRIPTION_* \
  --env ASTERISK_SIP_PORT \
  --env BRAND_APPID \
  --tz=${TIMEZONE} \
  ${NETHVOICE_FLEXISIP_IMAGE}
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/flexisip.ctr-id -t 10
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/flexisip.ctr-id
PIDFile=%t/flexisip.pid
Type=forking

[Install]
WantedBy=default.target
