# imageroot/systemd/user/watcher.path
#
# Watch for files named <action>_<service>
#
[Path]
PathExistsGlob=%S/state/notify/*_*

[Install]
WantedBy=default.target
