# imageroot/systemd/user/watcher.path
[Path]
PathModified=%S/state/notify/restart

[Install]
WantedBy=default.target
